# deploy/supervisor_weblog.ini
[program:uwsgi]
environment=WEBLOG_CACHE_DIR=/tmp/weblogcache
command=uwsgi --ini-paste-logged /var/weblog/weblog.ini
autorestart=true
redirect_stderr=true
stopsignal=INT