# deploy/supervisor_weblog.ini
[program:uwsgi]
environment=WEBLOG_CACHE_DIR=/root/weblogcache
command=uwsgi --ini-paste-logged /opt/weblog.ini
autorestart=true
redirect_stderr=true
stopsignal=INT